# supervisord.conf
[program:runserver]
exclude=True

[program:autoreload]
exclude=True

[program:celerybeat]
exclude=True
